# deployment/docker/supervisor/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:django]
command=gunicorn ai_productivity_platform.wsgi:application --bind 0.0.0.0:8000 --workers 4 --timeout 120 --worker-class sync --max-requests 1000 --max-requests-jitter 100
directory=/app/backend
user=appuser
autostart=true
autorestart=true
stdout_logfile=/app/logs/django.log
stderr_logfile=/app/logs/django_error.log
environment=DJANGO_SETTINGS_MODULE="ai_productivity_platform.settings"

[inet_http_server]
port=9001
username=admin
password=admin123

[supervisorctl]
serverurl=http://localhost:9001
username=admin
password=admin123

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface